# files/travel/wg-proton.service.d/roost.conf
[Unit]
After=network-online.target
Wants=network-online.target

[Service]
# Kill-switch sanity: iptables -C writes "Bad rule" to stderr on miss; suppress
# that generic noise so our WARN log is the only signal worth reading.
ExecStartPost=/bin/bash -c 'iptables -C OUTPUT -m owner --uid-owner xray ! -o wg-proton ! -o lo -j REJECT 2>/dev/null || logger -t roost/wg-proton "WARN: IPv4 kill-switch rule missing after PostUp"'
ExecStartPost=/bin/bash -c 'ip6tables -C OUTPUT -m owner --uid-owner xray ! -o wg-proton ! -o lo -j REJECT 2>/dev/null || logger -t roost/wg-proton "WARN: IPv6 kill-switch rule missing after PostUp"'
