# files/travel/wg-proton.service.d/roost.conf
[Unit]
After=network-online.target
Wants=network-online.target

[Service]
# Kill-switch sanity: fail the unit if the REJECT rule is absent so cmd_vpn's
# `enable --now` branch catches it. `||` chained to logger would always exit 0
# because logger returns 0; the `if then else exit 1` form propagates the miss.
ExecStartPost=/bin/bash -c 'if ! iptables -C OUTPUT -m owner --uid-owner xray -j REJECT 2>/dev/null; then logger -t roost/wg-proton "FATAL: IPv4 kill-switch REJECT rule missing after PostUp"; exit 1; fi'
ExecStartPost=/bin/bash -c 'if ! ip6tables -C OUTPUT -m owner --uid-owner xray -j REJECT 2>/dev/null; then logger -t roost/wg-proton "FATAL: IPv6 kill-switch REJECT rule missing after PostUp"; exit 1; fi'
