# files/glances.service
[Unit]
Description=Glances system monitoring
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
User=$USERNAME
ExecStartPre=+/bin/bash -c 'for i in $(seq 1 30); do tailscale ip -4 >/dev/null 2>&1 && exit 0; sleep 2; done; echo "ERROR: Tailscale IP not ready after 60s" >&2; exit 1'
ExecStart=/bin/bash -c '/usr/bin/glances -w --bind "$(tailscale ip -4)" --port 61208 --disable-plugin cloud'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
