# files/ram-monitor.timer
[Unit]
Description=Check for high RAM processes every 30s

[Timer]
OnBootSec=30
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
