# files/laptop/clip-forward.service
[Unit]
Description=Clipboard forwarding daemon for remote Claude Code
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/clip-forward serve
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
