# files/laptop/drop-watch.service
[Unit]
Description=Watch ~/drop/ and sync to roost server
After=network-online.target tailscaled.service

[Service]
Type=simple
ExecStart=%h/.local/bin/drop-watch
Restart=always
RestartSec=60
Environment=ROOST_SERVER=${SERVER_NAME}
Environment=ROOST_USER=${USERNAME}
Environment=ROOST_DIR_NAME=${ROOST_DIR_NAME}

[Install]
WantedBy=default.target
