# files/laptop/ntfy-cf-ip-refresh@.service
[Unit]
Description=Notify ntfy on cf-ip-refresh failure for %i

[Service]
Type=oneshot
# %i is the failed unit's name (passed via OnFailure=ntfy-cf-ip-refresh@%n
# in cf-ip-refresh.service). NTFY_URL + NTFY_TOKEN come from the env file
# populated by install-cf-ip-refresh.sh (server's Tailscale IP + the hooks
# admin token from server's ~/services/.ntfy-token; ntfy's auth-default-
# access is deny-all so anonymous posts are rejected).
EnvironmentFile=/etc/cf-ip-refresh.env
ExecStart=/usr/bin/curl -fsS -X POST -H "Authorization: Bearer ${NTFY_TOKEN}" -H "Title: cf-ip-refresh FAILED" -H "Priority: high" -H "Tags: warning" --data "Unit %i failed; check 'journalctl -u %i' for details." ${NTFY_URL}
