# files/travel/proton-keepalive.service
[Unit]
Description=Proton WireGuard keepalive probe
ConditionPathExists=/sys/class/net/wg-proton
After=wg-quick@proton.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/proton-keepalive-check
